# samples/playbooks/templates/wlsadminserver.conf.j2
start on runlevel [345]
exec /bin/su - {{ fmw_user }} -- {{ domains_home }}/{{ domain_name }}/bin/startWebLogic.sh
